# demo-nodes/app-node/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:sample-service]
command=/bin/sh -c "java -javaagent:/opt/skywalking/skywalking-agent.jar -Dskywalking.agent.service_name=sample-service -Dskywalking.collector.backend_service=skywalking-oap:11800 -jar /opt/app/sample-service.jar"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:node-exporter]
command=/usr/local/bin/node_exporter --path.rootfs=/ --web.listen-address=:9100
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:agent]
command=/usr/local/bin/monitor-agent
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
